# nyxid-node/rootfs/etc/services.d/nyxid-node/finish
#!/command/with-contenv bashio
# ==============================================================================
# NyxID Node Add-on — Finish (cleanup on stop/crash)
# ==============================================================================

if [ "${1}" -ne 0 ]; then
    bashio::log.warning "NyxID node agent exited with code ${1} (signal ${2})"
fi
